# wireguard_client/rootfs/etc/services.d/wireguard_client/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: WireGuard Client
# Runs WireGuard Client
# ==============================================================================
declare interface

s6-svc -O /run/service/wireguard_client

bashio::log.info "Starting WireGuard Client..."

# Get the interface
interface="wg0"

# Run the WireGuard
exec wg-quick up "${interface}"
